# deploy/tepinhui-backend.service
[Unit]
Description=Tepinhui Backend Application
After=network.target
StartLimitIntervalSec=60
StartLimitBurst=3
OnFailure=tepinhui-backend-rollback.service

[Service]
Type=simple
User=tph
WorkingDirectory=/home/tph/versions/latest
Environment="JAVA_OPTS=-Xms256m -Xmx512m"
EnvironmentFile=-/home/tph/shared/app.env
ExecStart=/bin/bash -lc 'exec /usr/bin/java $JAVA_OPTS -jar /home/tph/versions/latest/app.jar'
Restart=on-failure
RestartSec=10
StandardOutput=append:/home/tph/shared/logs/tepinhui-backend/application.log
StandardError=append:/home/tph/shared/logs/tepinhui-backend/application.log

[Install]
WantedBy=multi-user.target
